# supervisor/supervisor.conf
[unix_http_server]
file=/tmp/supervisor.sock   ; socket file
chmod=0700

[supervisord]
logfile=/Users/emilliacantella/kuliah/skripsi/adminbot/supervisor/supervisord.log
pidfile=/tmp/supervisord.pid
childlogdir=/Users/emilliacantella/kuliah/skripsi/adminbot/supervisor/

[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface

[supervisorctl]
serverurl=unix:///tmp/supervisor.sock

[program:laravel-worker]
process_name=%(program_name)s_%(process_num)02d
command=php /Users/emilliacantella/kuliah/skripsi/adminbot/artisan queue:work --sleep=3 --tries=3
autostart=true
autorestart=true
stopasgroup=true
killasgroup=true
user=emilliacantella
numprocs=1
redirect_stderr=true
stdout_logfile=/Users/emilliacantella/kuliah/skripsi/adminbot/storage/logs/worker.log

[program:laravel-reverb]
process_name=%(program_name)s
command=php /Users/emilliacantella/kuliah/skripsi/adminbot/artisan reverb:start
autostart=true
autorestart=true
stopasgroup=true
killasgroup=true
user=emilliacantella
redirect_stderr=true
stdout_logfile=/Users/emilliacantella/kuliah/skripsi/adminbot/storage/logs/reverb.log

